# user_services/scanner-mcp.service
[Unit]
Description=Scanner MCP Server (Warm Whisper GPU)
After=network.target

[Service]
Type=simple
EnvironmentFile=/home/ned/Documents/neds_scanner_radio_full_pipeline_with_git/transcriber/.env
WorkingDirectory=/home/ned/Documents/neds_scanner_radio_full_pipeline_with_git/transcriber
ExecStart=/home/ned/venv/bin/python3 -u /home/ned/Documents/neds_scanner_radio_full_pipeline_with_git/transcriber/scanner_transcriber_mcp.py
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
